# pkg/systemd/task-timer-server.service
# Task Timer gateway.
#
# This is a SYSTEM unit, unlike task-timer-sync.service. The two are not the same
# kind of thing and it is worth being explicit about why:
#
#   task-timer-sync   runs as a person, reads that person's database, and syncs
#                     their own work. There is nothing for it to do as root.
#   task-timer-server runs as a service. It holds the Atlassian OAuth grants for
#                     every user, listens on a port, and is nobody's login session.
#
# It carries no credentials in this file. The two secrets are passed as systemd
# credentials, which land in a file mode 0400 owned by the service user rather
# than in the environment, where /proc/<pid>/environ and every child process and
# core dump would have them.

[Unit]
Description=Task Timer gateway
Documentation=file:/etc/task-timer-server/config.toml
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=task-timer-server
Group=task-timer-server

ExecStart=/usr/bin/task-timer-server

# The Atlassian client secret and the token-encryption key. Put them in
# /etc/task-timer-server/ as root-owned 0600 files; systemd copies them into the
# unit's private credential store and hands the service a path to read.
#
# The encryption key must survive restarts and must never be regenerated: it
# decrypts the stored Jira refresh tokens, and a new key silently invalidates
# every one of them, forcing every user to reconnect with no explanation.
LoadCredential=atlassian_client_secret:/etc/task-timer-server/atlassian_client_secret
LoadCredential=token_encryption_key:/etc/task-timer-server/token_encryption_key

Restart=on-failure
RestartSec=5s

# It talks to the internet and holds other people's OAuth grants, so it gets a
# real sandbox — unlike the desktop daemon, whose data directory is relocatable
# and would break under one.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
MemoryDenyWriteExecute=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# ProtectSystem=strict makes all of /usr and /etc read-only, so the one directory
# it legitimately writes has to be named.
ReadWritePaths=/var/lib/task-timer-server
StateDirectory=task-timer-server
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
